# systemd/peoplesoft.service
[Unit]
Description=PeopleSoft PIA
After=network-online.target oracledb.service
Wants=network-online.target oracledb.service

[Service]
Type=oneshot
User=psoft
Group=psoft
ExecStart=/app/psoft/peopletools-admin/scripts/start_all.sh -y
ExecStop=/app/psoft/peopletools-admin/stop_all.sh -y
RemainAfterExit=yes
TimeoutSec=180

[Install]
WantedBy=multi-user.target