# systemd/peoplesoft.service
[Unit]
Description=PeopleSoft PIA
After=network-online.target oracledb.service
Wants=network-online.target oracledb.service

[Service]
Type=forking
User=psoft
Group=psoft

ExecStart=/app/psoft/peopletools-admin/systemd/service_start.sh
ExecStop=/app/psoft/peopletools-admin/systemd/service_stop.sh

RemainAfterExit=yes

[Install]
WantedBy=multi-user.target